# src/usr/lib/systemd/system/factorio-manager-reset@.service
[Unit]
Description="Factorio Server reset (%i)"

[Service]
ExecStart=/usr/bin/systemctl kill -s USR1 factorio-manager@%i.service
